# crontab
8 */6 * * * /usr/bin/python /path/to/hack0verflow.py
